# pocketchip-batt/system/pocketchip-warn15.service
[Unit]
Description=15% battery warning polling

[Service]
Type=oneshot
User=chip
ExecStart=/usr/sbin/pocketchip-warn15
