# pocketchip-batt/system/pocketchip-batt-warn.timer
[Unit]
Description=poll battery for warning every 30 seconds

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
