# pocketchip-batt/system/pocketchip-warn-clean.service
[Unit]
Requires=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'rm /usr/lib/pocketchip-batt/warn*'
