# pocketchip-batt/system/pocketchip-batt.timer
[Unit]
Description=poll battery and charging status every 3 seconds

[Timer]
OnBootSec=10sec
OnUnitActiveSec=2sec

[Install]
WantedBy=timers.target
