# pocketchip-batt/system/pocketchip-warn05.service
[Unit]
Description=5% battery warning polling

[Service]
Type=oneshot
ExecStart=/usr/sbin/pocketchip-warn05
