# pocketchip-batt/system/pocketchip-batt.service
[Unit]
Description=battery polling for pocket-wm

[Service]
Type=oneshot
ExecStart=/usr/sbin/pocketchip-batt
