# pocketchip-batt/system/pocketchip-off00.timer
[Unit]
Description=poll battery for 00% shutdown

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
