# pocketchip-batt/system/pocketchip-warn15.timer
[Unit]
Description=poll battery for 15% warning every 30 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
