# pocketchip-batt/system/pocketchip-off00.service
[Unit]
Description=00% battery warning polling

[Service]
Type=oneshot
User=chip
ExecStart=/usr/sbin/pocketchip-off00
