# pocketchip-batt/system/pocketchip-warn05.timer
[Unit]
Description=poll battery for 5% warning every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
